# files/default/sample_app.conf
description "Sample app service for test"

respawn
respawn limit 15 5

start on runlevel [2345]
stop on runlevel [016]

setuid app

exec /usr/local/bin/sample_app
